# scripts/mykey-tray.service
[Unit]
Description=MyKey System Tray
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/mykey-tray run
Restart=on-failure
RestartSec=3s

[Install]
WantedBy=graphical-session.target
